# install/ddmgr.service
[Unit]
Description=Generate network configs
PartOf=isc-dhcp-server.service named.service nftables.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ddmgr_config

[Install]
RequiredBy=isc-dhcp-server.service named.service nftables.service
